# services/mediamtx.service
[Unit]
Description=MediaMTX RTSP/RTMP/SRT Server
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/mediamtx
ExecStart=/usr/local/bin/mediamtx /opt/mediamtx/mediamtx.yml
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

