# services/preke-recorder.service
[Unit]
Description=Preke R58 Recorder Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/preke-r58-recorder
Environment="PATH=/opt/preke-r58-recorder/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONPATH=/opt/preke-r58-recorder"
# Cloudflare TURN credentials for remote WebRTC guests
Environment="CLOUDFLARE_TURN_TOKEN_ID=79d61c83455a63d11a18c17bedb53d3f"
Environment="CLOUDFLARE_TURN_API_TOKEN=9054653545421be55e42219295b74b1036d261e1c0259c2cf410fb9d8a372984"
# Disable GStreamer plugin scanning at startup (use cached registry)
Environment="GST_REGISTRY_UPDATE=no"
# Clean up stuck GStreamer processes before starting (use - prefix to ignore failures)
ExecStartPre=-/usr/bin/pkill -9 gst-plugin-scanner
ExecStartPre=-/bin/sleep 1
# Ensure recordings directory symlink exists (config uses /mnt/sdcard/recordings, NVMe SSD is at /data)
ExecStartPre=-/bin/mkdir -p /data/recordings
ExecStartPre=-/bin/ln -sfn /data/recordings /mnt/sdcard/recordings
# Initialize HDMI inputs (query subdev resolution and set format on video devices)
ExecStartPre=-/opt/preke-r58-recorder/scripts/init_hdmi_inputs.sh
ExecStart=/opt/preke-r58-recorder/venv/bin/uvicorn src.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=10
# Timeout for startup - if service doesn't start in 60s, consider it failed
TimeoutStartSec=60
# Shorter stop timeout to avoid long waits
TimeoutStopSec=30
# Kill entire process group on stop
KillMode=control-group
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
