# services/r58-pipeline.service
[Unit]
Description=R58 Pipeline Manager
After=network.target

[Service]
Type=simple
User=r58
Group=r58
WorkingDirectory=/opt/r58
ExecStart=/opt/r58/venv/bin/python -m pipeline_manager.main
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
RuntimeDirectory=r58
Environment=PYTHONPATH=/opt/r58/packages/backend

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/r58 /opt/r58/recordings /run/r58

[Install]
WantedBy=multi-user.target

