# systemd/r58-pipeline.service
[Unit]
Description=R58 Pipeline Manager
After=network.target
Before=r58-api.service

[Service]
Type=simple
User=linaro
# IMPORTANT: Use /opt/preke-r58-recorder directly, NOT the symlink
# The symlink /opt/r58-app/current caused deployment issues where old code kept running
WorkingDirectory=/opt/preke-r58-recorder
ExecStart=/opt/preke-r58-recorder/venv/bin/python -m pipeline_manager.main
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal

# Create /run/r58 directory on startup
RuntimeDirectory=r58
RuntimeDirectoryMode=0755

# Environment
Environment="PATH=/opt/preke-r58-recorder/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONPATH=/opt/preke-r58-recorder/packages/backend"
Environment="PYTHONUNBUFFERED=1"

# Security (relaxed for linaro user)
ReadWritePaths=/opt/preke-r58-recorder/recordings /run/r58

[Install]
WantedBy=multi-user.target
