# systemd/r58-pipeline.service
[Unit]
Description=R58 Pipeline Manager
After=network.target
Before=r58-api.service

[Service]
Type=simple
User=linaro
WorkingDirectory=/opt/r58-app/current
ExecStart=/opt/r58-app/current/venv/bin/python -m pipeline_manager.main
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal

# Create /run/r58 directory on startup
RuntimeDirectory=r58
RuntimeDirectoryMode=0755

# Environment
Environment="PATH=/opt/r58-app/current/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONPATH=/opt/r58-app/current/packages/backend"
Environment="PYTHONUNBUFFERED=1"

# Security (relaxed for linaro user)
ReadWritePaths=/opt/r58-app/current/recordings /run/r58

[Install]
WantedBy=multi-user.target
