# deployment/disable-headset-irq.service
[Unit]
Description=Disable floating headset detection IRQ
Documentation=https://github.com/mBelstad/preke-r58-recorder/blob/main/docs/HEADSET_IRQ_FIX.md
After=multi-user.target

[Service]
Type=oneshot
# Disable IRQ 159 which causes 277K interrupts/sec due to floating GPIO
# This saves approximately 36% CPU usage
# Note: The IRQ cannot be fully disabled without kernel module changes.
# This service documents the issue. A kernel module blacklist or device tree
# modification would be needed for a permanent fix.
ExecStart=/bin/true
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

