# deployment/vdo-nginx-webapp.service
[Unit]
Description=VDO.ninja Web App Server
After=network.target

[Service]
Type=simple
User=linaro
WorkingDirectory=/opt/vdo.ninja
ExecStart=/usr/bin/python3 -m http.server 8444 --bind 0.0.0.0
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

