# scripts/vdoninja-bridge.service
[Unit]
Description=VDO.ninja WHEP Bridge - Share HDMI cameras to VDO.ninja room (background)
After=network-online.target graphical.target preke-recorder.service
Wants=network-online.target
Requires=graphical.target
# Start before TV kiosk so kiosk appears on top
Before=preke-tv-kiosk.service

[Service]
Type=oneshot
RemainAfterExit=yes
User=linaro
Group=linaro

# Display configuration
Environment=DISPLAY=:0

# VDO.ninja room name
# Must match VDO_ROOM in packages/frontend/src/lib/vdoninja.ts
Environment=VDONINJA_ROOM=studio

# VDO.ninja host (self-hosted instance, now on same domain as MediaMTX)
Environment=VDONINJA_HOST=app.itagenten.no/vdo

# R58 API host (for WHEP streams)
Environment=API_HOST=app.itagenten.no

# Camera configuration - comma-separated list
# Format: stream_id:push_id:label
# - stream_id: MediaMTX stream name (cam0, cam1, cam2, cam3)
# - push_id: VDO.ninja push ID (must be unique per camera)
# - label: Display name shown in VDO.ninja
#
# Example for single camera:
#   CAMERAS=cam2:hdmi1:HDMI-Camera
#
# Example for multiple cameras:
#   CAMERAS=cam2:hdmi1:Camera-1,cam3:hdmi2:Camera-2
#
# Default: All 3 active HDMI cameras
Environment="CAMERAS=cam0:hdmi0:HDMI-IN0,cam2:hdmi2:HDMI-IN1,cam3:hdmi3:HDMI-IN2"

# Program Output configuration
# When enabled, the device renders the VDO.ninja mixed scene and publishes
# to MediaMTX, enabling RTMP relay to YouTube/Twitch
Environment="ENABLE_PROGRAM_OUTPUT=true"
Environment="PROGRAM_OUTPUT_PATH=mixer_program"
Environment="MEDIAMTX_PUBLISH_HOST=app.itagenten.no:443"

# Iframe sources (auto-added to VDO.ninja room)
# Format: "push_id|label|url" (pipe-delimited to avoid URL colon conflicts)
# These are embedded as iframe video sources and pushed to the room automatically
# Examples:
#   slides|Slides|https://app.itagenten.no/reveal - Reveal.js presentation
#   graphics|Graphics|https://example.com/graphics - CasparCG-style graphics
Environment="IFRAME_SOURCES=slides|Slides|https://app.itagenten.no/reveal"

# Log file location
Environment=LOG_FILE=/var/log/vdoninja-bridge.log

WorkingDirectory=/opt/preke-r58-recorder
# Reduced from 10s to 3s - display/network usually ready faster
ExecStartPre=/bin/sleep 3
ExecStart=/opt/preke-r58-recorder/scripts/start-vdoninja-bridge.sh
StandardOutput=journal
StandardError=journal

# NOTE: This service should NOT be enabled at boot!
# The mode_manager.py controls starting/stopping this service
# based on the current mode (starts in mixer mode, stops in recorder/idle)
[Install]
WantedBy=multi-user.target
